# systemd-template.service
[Unit]
Description="--Description--"

[Service]
Type=oneshot|simple
ExecStart=/usr/bin/script

[Install]
WantedBy=multi-user.target
